# sports_live/rootfs/etc/services.d/sports-live/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Start the Sports Live add-on
# ==============================================================================
bashio::log.info "Starting Sports Live add-on..."

cd /opt/app
exec python -m app
